# postgresql14/rootfs/etc/systemd/system/postgresql-14.service.d/override.conf
[Unit]
After=postgresql-initdb.service
Requires=postgresql-initdb.service
