# rabbitmq/rootfs/etc/systemd/system/rabbitmq-server.service.d/override.conf
[Unit]
Before=rabbitmq-config.service
Requires=rabbitmq-config.service
