# nginx/rootfs/etc/systemd/system/etc-nginx-conf.d.path
[Unit]
Description=Monitor nginx conf.d changes

[Path]
PathChanged=/etc/nginx/conf.d
Unit=nginx-reload.service

[Install]
WantedBy=multi-user.target
