# nginx/rootfs/etc/systemd/system/etc-nginx-stream.d.path
[Unit]
Description=Monitor nginx stream.d changes

[Path]
PathChanged=/etc/nginx/stream.d
Unit=nginx-reload.service

[Install]
WantedBy=multi-user.target
