# postgresql14/rootfs/etc/systemd/system/postgresql-password.service
[Unit]
Description=PostgreSQL 14 password generator
ConditionFileNotEmpty=!/run/secrets/postgres-password

[Service]
Type=oneshot
ExecStart=/usr/bin/bash -c "printf \"$(cat /proc/sys/kernel/random/uuid)\n\" > /run/secrets/postgres-password"
RemainAfterExit=false
