# grafana/rootfs/etc/systemd/systemd/grafana-server.service.d/override.conf
[Unit]
After=systemd-tmpfiles-setup.service
Requires=systemd-tmpfiles-setup.service
