# postgresql13/rootfs/etc/systemd/system/postgresql-13.service.d/override.conf
[Unit]
After=postgresql-config.service
Requires=postgresql-config.service
