# couchdb/rootfs/etc/systemd/system/couchdb.service.d/override.conf
[Unit]
Requires=couchdb-password.service couchdb-config.service
After=couchdb-password.service couchdb-config.service
