# mysql8/rootfs/etc/systemd/system/mysqld.service.d/override.conf
[Unit]
Before=mysql-config.service
Requires=mysql-config.service
